# examples/echo.service
[Service]
ExecStart=/home/carl/code/systemd.cr/echo
Type=notify
FileDescriptorStoreMax=infinity
